# systemdunits/battery_monitor.timer
[Unit]
Description=Battery Monitor Timer

[Timer]
OnCalendar=minutely
Unit=battery_monitor.service

[Install]
WantedBy=timers.target
