# baseconfig/systemd/.config/systemd/user/battery_monitor.service
[Unit]
Description=Battery Monitor Service
[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery_warning.py
[Install]
WantedBy=multi-user.target
